# cleanup-samba-recycle.service.example
[Unit]
Description=Cleanup Samba recycle bin

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cleanup-samba-recycle.sh
